# foodshop/backend/deploy/foodshop-backend.service
[Unit]
Description=Foodshop Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/foodshop/backend
Environment=DB_PATH=/opt/foodshop/backend/data/foodshop.db
Environment=DEVICE_TOKENS=kitchen-tablet-001
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
